# deploy/aws/quantchat-backend.service
[Unit]
Description=QuantChat FastAPI Backend
After=network.target

[Service]
User=ubuntu
WorkingDirectory=__APP_DIR__/backend
EnvironmentFile=__APP_DIR__/backend/.env
ExecStart=__APP_DIR__/backend/.venv/bin/uvicorn server:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
